# services/extra/crontab
* * * * * /usr/local/bin/php /www/wwwroot/artisan schedule:run